# templates/default/jmxtrans.init.el.erb
#!/bin/sh
#
# chkconfig: - 90 60
# pidfile: /var/run/<%= @name %>.pid
# config: /etc/default/<%= @name %>
#
### BEGIN INIT INFO
# Provides: <%= @name %>
# Required-Start: $local_fs $remote_fs $network $named
# Required-Stop: $local_fs $remote_fs $network
# Short-Description: run <%= @name %> daemon
# Description: <%= @name %> is a very powerful tool which reads json
#              configuration files of servers/ports and jmx domains - attributes - types
#              and then outputs the data in whatever format you want via special 'Writer' objects which you can code up yourself.
#              It does this with a very efficient engine design that will scale to querying thousands of machines.
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions
NAME=<%= @name %>

if [ -f /etc/default/$NAME ]; then
    . /etc/default/$NAME
fi

RETVAL=0
prog=$NAME
user=$NAME
exec='<%="#{node['jmxtrans']['home']}/jmxtrans.sh" %>'
lockfile=/var/lock/subsys/<%= @name %>
config=/etc/default/<%= @name %>

start() {
    echo -n $"Starting $prog: "
    daemon --user=$user $exec start
    #$exec start
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch $lockfile
}

stop() {
    echo -n $"Stopping $prog: "
    $exec stop
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f $lockfile
}

status() {
    $exec status
}

restart() {
    daemon --user=$user $exec restart
}

# See how we were called.
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        RETVAL=$?
        ;;
  restart|force-reload|reload)
        restart
        ;;
  *)
        echo $"Usage: $prog {start|stop|restart|force-reload|reload|status}"
        RETVAL=2
esac

exit $RETVAL
